# misc/crontab
# m h  dom mon dow   command
0 * * * * cd /root && /usr/local/bin/r8 sql stmt SELECT 1 >> /var/log/r8-backup.log 2>&1
